# hearthloop.service
[Unit]
Description=HearthLoop Dashboard Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/hearthloop
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
